# supervisord.conf
[supervisord]
nodaemon=true

[program:sockexec]
command=/opt/bin/sockexec -t0 /tmp/exec.sock

